# flexvdi-guest-agent.service.in
[Unit]
Description=flexVDI Guest Agent

[Service]
ExecStartPre=/bin/rm -f /var/run/flexvdi_pipe /var/run/flexvdi-guest-agent.pid
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/flexvdi-guest-agent
Type=forking
PIDFile=/var/run/flexvdi-guest-agent.pid
Restart=on-abort

[Install]
WantedBy=multi-user.target
